# bin/crontab/crontab
MAILTO=admin@domain.tld
* * * * * php -v
MAILTO=admin@domain.tld
* * * * * php -i
